# matrix-display.service
[Unit]
Description=LED Matrix Display Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=INSTALL_DIR
ExecStart=/usr/bin/python3 INSTALL_DIR/server.py
# Konfiguration via config.json – keine Env-Vars nötig
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
